# roles/caddy/templates/caddy.service.j2
[Unit]
Description=Caddy
Documentation=https://caddyserver.com/docs/
After=network.target network-online.target
Requires=network-online.target

[Service]
Type=notify
User={{ caddy_user }}
Group={{ caddy_group }}
EnvironmentFile={{ caddy_config_dir }}/caddy.env
ExecStart={{ caddy_install_path }} run --environ --config {{ caddy_config_dir }}/Caddyfile
ExecReload={{ caddy_install_path }} reload --config {{ caddy_config_dir }}/Caddyfile --force
TimeoutStopSec=5s
LimitNOFILE=1048576
PrivateTmp=true
ProtectSystem=full
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
